# debian_files/wdb-1.6.wdb.cron.d
SHELL=/bin/sh


36 1 * * * wdb /usr/lib/wdb/wdb/clean_wdb.sh

