# debian_files/wdb.cron.d
SHELL=/bin/sh


15 4 * * * su postgres -c "/usr/lib/wdb/clean_wdb.sh"
