# debian_files/wdb.cron.d
SHELL=/bin/sh


36 2 * * * wdb /usr/lib/wdb/clean_wdb.sh

